# yt-dlp/etc/services.d/yt-dlp/run
#!/usr/bin/with-contenv bash

cd /app
umask "$UMASK"
exec s6-setuidgid abc npm start